# ansible-home/.config/systemd/user/ansible.service
[Unit]
Description=Run Ansible playbook
After=network.target

[Service]
User=your_username
WorkingDirectory=~/.ansible
ExecStart=ansible-playbook -i hosts deploy.yml

