# ansible-home/.config/systemd/user/ansible.timer
[Unit]
Description=Run Ansible playbook every 5 minutes

[Timer]
OnUnitActiveSec=5min
Unit=ansible.service

[Install]
WantedBy=timers.target

